# bin/assets/aspiesoft-clamav-scanner/aspiesoft-clamav-download-scanner.service
sudo mkdir -p "$HOME/.local/share/Steam/steamapps"[Unit]
Description=AspieSoft ClamAV Download Scanner

[Service]
ExecStart=/bin/bash /etc/aspiesoft-clamav-scanner/init.sh

[Install]
WantedBy=multi-user.target
